# systemd/unet-server.service
[Unit]
Description=μNet Network Configuration Management Server
Documentation=https://github.com/example/unet
After=network.target network-online.target
Wants=network-online.target
RequiresMountsFor=/var/lib/unet

[Service]
Type=exec
User=unet
Group=unet
ExecStart=/usr/bin/unet-server --config /etc/unet/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
CapabilityBoundingSet=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@debug @mount @cpu-emulation @obsolete @privileged @reboot @swap @raw-io

# Allow writing to configuration and data directories
ReadWritePaths=/var/lib/unet /var/log/unet /etc/unet

# Working directory
WorkingDirectory=/var/lib/unet

# Environment
Environment=RUST_LOG=info
Environment=UNET_CONFIG_DIR=/etc/unet
Environment=UNET_DATA_DIR=/var/lib/unet
Environment=UNET_LOG_DIR=/var/log/unet

# Limits
LimitNOFILE=65536
LimitNPROC=4096
LimitCORE=0

# Watchdog
WatchdogSec=30

# OOM protection
OOMScoreAdjust=-900

[Install]
WantedBy=multi-user.target